# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:rabbitmq]
command=rabbitmq-server
